# unifi2mqtt/rootfs/etc/services.d/unifi2mqtt/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting unifi2mqtt..."
exec pm2-runtime start -- runuser --preserve-environment -u unifi2mqtt unifi2mqtt
