# json2mqtt/rootfs/etc/services.d/json2mqtt/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting json2mqtt..."
exec runuser --preserve-environment -u json2mqtt /run.sh
