# kibana/rootfs/etc/services.d/kibana/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e +u

source "${ENVFILE:-/usr/share/kibana/.env}"
exec /usr/local/bin/kibana-docker
