# kibana/rootfs/etc/services.d/kibana/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e +u

export KBN_HOME=/usr/share/kibana
export KBN_PATH_CONF=/etc/kibana
export RESTART_ON_UPGRADE="true"

cd "${KBN_HOME}" || exit 1

s6-setuidgid kibana /usr/share/kibana/bin/kibana
